# ops/global-documentation-sync/global-documentation-sync.timer
[Unit]
Description=Run QueryPie Global publication sync daily

[Timer]
OnCalendar=*-*-* 10:00:00
Persistent=true
RandomizedDelaySec=10m
Unit=global-documentation-sync.service

[Install]
WantedBy=timers.target
