# ops/global-documentation-sync/global-documentation-sync.timer
[Unit]
Description=Run QueryPie Global publication sync on weekdays

[Timer]
OnCalendar=Mon..Fri *-*-* 10:00:00 Asia/Seoul
Persistent=true
RandomizedDelaySec=10m
Unit=global-documentation-sync.service

[Install]
WantedBy=timers.target
